# zoom-update-sangha/etc/systemd/system/zoom_update_sangha.timer
[Unit]
Description=Run Zoom update script daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
